# templates/ddos.service.j2
[Unit]
After=network.target network-online.target

[Service]
ExecStart={{ ddos }}
User={{ ansible_user }}

[Install]
WantedBy=default.target


